# data-restore/backup.cron
0 0 * * 1-6 root /backup.sh incremental
0 0 * * 7 root /backup.sh full
#