# data-backup/backup.cron
0 0 1 * * root /backup.sh full
0 0 2-31 * * root /backup.sh incremental
#